# supervisor.ini
[program:blog]
command=/data/env-blog/bin/uwsgi /data/www/blog/config.ini
directory=/data/www/blog
autostart=true
autorestart=true
stdout_logfile=/data/logs/uwsgi_supervisor.log